# packages/pglift/pglift-0.22.0.tar.gz/pglift-0.22.0/src/pglift/data/systemd/pglift-postgres_exporter@.service
[Unit]
Description=Prometheus exporter %i database server metrics
Documentation=https://github.com/wrouesnel/postgres_exporter
StartLimitIntervalSec=60

[Service]
Type=simple
{executeas}

PrivateTmp=yes
NoNewPrivileges=yes

EnvironmentFile=-{configpath}

ExecStart={execpath} $POSTGRES_EXPORTER_OPTS
ExecReload=/bin/kill --signal HUP $MAINPID
KillMode=process
KillSignal=SIGINT

Restart=on-failure
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
